# blademedia-organizer.service
[Unit]
Description=BladeMedia Media Organizer
After=docker.service network-online.target
Wants=docker.service

[Service]
Type=oneshot
WorkingDirectory=/srv/blademedia
ExecStart=/usr/bin/python3 /srv/blademedia/media_organizer.py . --discord-webhook="YOUR_DISCORD_WEBHOOK_URL_HERE"
User=1000
Group=1000
StandardOutput=append:/srv/blademedia/organizer.log
StandardError=append:/srv/blademedia/organizer.log
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
